# .config/crontab
* * * * * kill -39 $(pidof dwmblocks)
0 0 * * * kill -38 $(pidof dwmblocks)
15 0 * * * notify-send 'Vai dormir muleke!!!!'
