# .config/crontab
* * * * * kill -39 $(pidof dwmblocks)
2 0 0 * * * kill -38 $(pidof dwmblocks)
3 0 0 * * * notify-send 'Vai dormir muleke!!!!'
